# postfix/assets/supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=warn

[program:postfix]
command=/usr/sbin/postfix -v start-fg

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
